# docker/s6/services.d/app/run
#!/usr/bin/execlineb -P

with-contenv
#todo signals
php /react/app.php
